# clash.ini
[program:clash]
directory=/tmp
command=clash               
autostart=true
autorestart=true
startsecs=3
user = root
redirect_stderr = true
stdout_logfile=/tmp/clash.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 20
